# Packaging.Linux/coder-desktop.service
[Unit]
Description=Coder Desktop VPN Service
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/lib/coder-desktop/CoderVpnService
Restart=on-failure
RestartSec=10
RuntimeDirectory=coder-desktop
# Needs root for TUN device creation
User=root

[Install]
WantedBy=multi-user.target
